# backend/deploy/zhejian-api.service
[Unit]
Description=Zhejian API (辙见 · geo.simplewin.cn)
After=network.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/zhejian/backend
EnvironmentFile=/var/www/zhejian/backend/.env
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
